# templates/dev.cassandra.plist.erb
<?xml version="1.0" encoding="UTF-8"?>

<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>dev.cassandra</string>

    <key>Boxen</key>
    <dict>
      <key>Description</key>
      <string>Distributed database action.</string>

      <key>Ports</key>
      <array>
        <integer><%= scope.lookupvar "cassandra::config::storage_port" %></integer>
        <integer><%= scope.lookupvar "cassandra::config::ssl_storage_port" %></integer>
        <integer><%= scope.lookupvar "cassandra::config::rpc_port" %></integer>
      </array>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string><%= scope.lookupvar "cassandra::config::executable" %></string>
      <string>-f</string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <true />

    <key>UserName</key>
    <string><%= scope.lookupvar "::boxen_user" %></string>

    <key>WorkingDirectory</key>
    <string><%= scope.lookupvar "cassandra::config::datadir" %></string>

    <key>StandardErrorPath</key>
    <string><%= scope.lookupvar "cassandra::config::logdir" %>/stderr.log</string>

    <key>StandardOutPath</key>
    <string><%= scope.lookupvar "cassandra::config::logdir" %>/stdout.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>CASSANDRA_CONF</key>
      <string><%= scope.lookupvar "cassandra::config::configdir" %></string>
    </dict>

  </dict>
</plist>
